# templates/screen_on.timer.j2
[Unit]
Description=Standby screen out of office hours

[Timer]
OnCalendar=*-*-* {{ screen_switch_on_time }}
Unit=screen_switch@on.service

[Install]
WantedBy=multi-user.target
