# files/screen_switch@.service
[Unit]
Description=Switch screen state to %I

[Service]
Type=oneshot
ExecStart=/usr/local/bin/screen_switch.sh %I
